# checkbox-ng/plainbox/data/plainbox-ci-mailer.conf
# Plainbox mini CI mailer.
#
# This waits for a "plainbox-finished" event and then
# uses curl to post one of two possible log files
# to a configurable CGI URL. The CGI can do whatever it
# wants, such as emailing the report.
#
# An /etc/default/plainbox-ci-mailer config file
# with the SUBMIT_CGI variable defined is expected.

description     "Plainbox mini ci mailer"

start on plainbox-finished

task

env PLAINBOX_UPSTART_LOG=/var/log/upstart/plainbox.log
env PLAINBOX_DESKTOP_LOG=/tmp/plainbox.log

pre-start script
    [ -f /etc/default/plainbox-ci-mailer ] || exit 1
    [ -x /usr/bin/curl ] || exit 1
end script

script
    . /etc/default/plainbox-ci-mailer
    [ -z "$SUBMIT_CGI" ] && exit 1
    RELEASE=$(lsb_release -ds)
    if [ -f $PLAINBOX_UPSTART_LOG ]; then
        MESSAGE=$PLAINBOX_UPSTART_LOG
        SUBJECT="Plainbox CI testing run for $RELEASE server"
   elif [ -f $PLAINBOX_DESKTOP_LOG ]; then
        MESSAGE=$PLAINBOX_DESKTOP_LOG
        SUBJECT="Plainbox CI testing run for $RELEASE desktop"
    else
        MESSAGE="Something failed and plainbox didn't even start."
        SUBJECT="FAILED plainbox CI testing run for $RELEASE"
    fi
    if [ -f "$MESSAGE" ] ; then
        curl -F subject="$SUBJECT" -F plainbox_output=@$MESSAGE $SUBMIT_CGI
    else
        curl -F subject="$SUBJECT" -F plainbox_output="$MESSAGE" $SUBMIT_CGI
    fi
end script
