# checkbox-ng/plainbox/data/plainbox.conf
# plainbox sru - Plainbox
#
# Run SRU testing with Plainbox 60 seconds after networking
# has started to wait for DHCP to complete.

description     "Plainbox SRU"

# Make sure we start after networking
start on runlevel [2345]

stop on runlevel [!023456]

task

pre-start script
  plainbox check-config || exit 1
  mkdir -p /var/cache/plainbox
end script

exec /usr/bin/plainbox sru
